# snapshots/SNAPSHOT_20260628T124655Z_pre_phase3_phase5_storagebox_watchdog/infra/systemd/v2-telegram-listener.service
# V2 Telegram Command Listener — Phase 3 Safeguard D
# ====================================================
#
# Manual override listener allowing Sebastien to emit /v2_flat, /v2_resume,
# and /v2_status commands from Telegram (iPhone) without needing SSH.
#
# Deployment on VPS Hetzner (5.161.246.190)
# -----------------------------------------
#   1. Copy this file to /etc/systemd/system/v2-telegram-listener.service
#        sudo cp /home/badoun/trading-bot-v2/infra/systemd/v2-telegram-listener.service \
#                /etc/systemd/system/
#   2. Reload systemd:
#        sudo systemctl daemon-reload
#   3. Enable on boot + start:
#        sudo systemctl enable v2-telegram-listener.service
#        sudo systemctl start v2-telegram-listener.service
#   4. Verify running:
#        sudo systemctl status v2-telegram-listener.service
#        tail -F /home/badoun/trading-bot-v2/live/logs/telegram_listener.log
#   5. From Telegram, send /v2_help to confirm the listener replies.
#
# Pre-requisites
# --------------
#   - /home/badoun/.config/badoun/telegram.env present with TELEGRAM_BOT_TOKEN
#     and TELEGRAM_CHAT_ID (or V2_TG_BOT_TOKEN / V2_TG_CHAT_ID for isolation)
#   - chmod 600 on telegram.env (preserve from VPS deployment standard)
#   - Python 3.9+ with `requests` library installed system-wide
#
# Stop / restart
# --------------
#   sudo systemctl stop v2-telegram-listener.service
#   sudo systemctl restart v2-telegram-listener.service
#
# Uninstall
# ---------
#   sudo systemctl disable --now v2-telegram-listener.service
#   sudo rm /etc/systemd/system/v2-telegram-listener.service
#   sudo systemctl daemon-reload

[Unit]
Description=V2 Telegram Command Listener (Phase 3 Safeguard D)
Documentation=https://github.com/Badounofficial/trading-bot-v2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=badoun
Group=badoun
WorkingDirectory=/home/badoun/trading-bot-v2

# Load Telegram credentials (TELEGRAM_BOT_TOKEN, TELEGRAM_CHAT_ID,
# and/or V2_TG_BOT_TOKEN, V2_TG_CHAT_ID for V2-isolated config).
# The Python loader also reads this file directly (defensive double-load),
# so systemd EnvironmentFile + paper_trading.config _load_env both work.
EnvironmentFile=/home/badoun/.config/badoun/telegram.env

# Ensure project root is on PYTHONPATH so `from paper_trading import ...`
# works when invoked outside the project root via systemd.
Environment=PYTHONPATH=/home/badoun/trading-bot-v2
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 /home/badoun/trading-bot-v2/live/telegram_command_listener.py

# Resilience: restart on any non-clean exit, but back off so a crash loop
# doesn't hammer Telegram API.
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=10

# Hardening (minimal — listener only needs network + state dir writes).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/badoun/trading-bot-v2/live/state /home/badoun/trading-bot-v2/live/logs
ProtectHome=read-only

# Logging: stdout/stderr go to the listener's own log file (the listener
# also writes structured logs via Python logging).
StandardOutput=append:/home/badoun/trading-bot-v2/live/logs/telegram_listener_stdout.log
StandardError=append:/home/badoun/trading-bot-v2/live/logs/telegram_listener_stderr.log

[Install]
WantedBy=multi-user.target
