# snapshots/SNAPSHOT_20260628T124935Z_post_phase3_phase5_storagebox_watchdog/infra/systemd/v2-daily-reconciliation.service
# V2 Daily Reconciliation — Phase 3 Safeguard C
# ===============================================
# One-shot service triggered by v2-daily-reconciliation.timer at 12:05 UTC.
# Reads daemon state, computes 24h delta + DD + deviation, sends Telegram.
#
# Deployment alongside the timer file:
#   sudo cp /home/badoun/trading-bot-v2/infra/systemd/v2-daily-reconciliation.service \
#           /home/badoun/trading-bot-v2/infra/systemd/v2-daily-reconciliation.timer \
#           /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now v2-daily-reconciliation.timer
#   systemctl list-timers v2-daily-reconciliation
#   journalctl -u v2-daily-reconciliation.service --since today

[Unit]
Description=V2 Daily Reconciliation Report (Phase 3 Safeguard C)
Documentation=https://github.com/Badounofficial/trading-bot-v2
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=badoun
Group=badoun
WorkingDirectory=/home/badoun/trading-bot-v2
EnvironmentFile=/home/badoun/.config/badoun/telegram.env
Environment=PYTHONPATH=/home/badoun/trading-bot-v2
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 /home/badoun/trading-bot-v2/live/daily_reconciliation.py

# Hardening — script only needs read of state + write to its own history file
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/badoun/trading-bot-v2/live/state /home/badoun/trading-bot-v2/live/logs
ProtectHome=read-only

StandardOutput=append:/home/badoun/trading-bot-v2/live/logs/daily_reconciliation.log
StandardError=append:/home/badoun/trading-bot-v2/live/logs/daily_reconciliation.log

# Don't restart — this is oneshot. The timer fires daily anyway.

[Install]
WantedBy=multi-user.target
