# snapshots/SNAPSHOT_20260628T024854Z_pre_phase3_phase4_saturday_recap_G/infra/systemd/v2-daily-reconciliation.timer
# V2 Daily Reconciliation Timer — Phase 3 Safeguard C
# =====================================================
# Triggers v2-daily-reconciliation.service every day at 12:05 UTC.
#
# The 12:05 UTC slot is chosen specifically to be:
#   - just after the daemon's daily heartbeat (12:00 UTC) so the daily
#     reconciliation can verify the heartbeat fired and reconcile against
#     the same equity snapshot
#   - inside Hyperliquid's high-liquidity hours
#   - independent of funding-payment hours (HL pays hourly, so no special slot)
#
# Persistent=true ensures that if the VPS is down at 12:05, the timer fires
# as soon as the system is back up (catch-up behavior).

[Unit]
Description=V2 Daily Reconciliation Timer — fires at 12:05 UTC daily
Documentation=https://github.com/Badounofficial/trading-bot-v2
Requires=v2-daily-reconciliation.service

[Timer]
OnCalendar=*-*-* 12:05:00 UTC
Persistent=true
Unit=v2-daily-reconciliation.service
AccuracySec=1min

[Install]
WantedBy=timers.target
